# deploy/hypo.service
[Unit]
Description=Mortgage Refinance Calculator (FastAPI)
After=network.target

[Service]
User=flaim
WorkingDirectory=/www/hypo
EnvironmentFile=/www/hypo/.env
ExecStart=/www/hypo/.venv/bin/uvicorn api.app:app \
    --host 127.0.0.1 \
    --port 8501
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
